# com.maekawa.video-direction-audit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.maekawa.video-direction-audit</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/maekawanaoto/AIťĖčÁôļ10/venv/bin/python3</string>
        <string>/Users/maekawanaoto/AIťĖčÁôļ10/scripts/run_audit.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/maekawanaoto/AIťĖčÁôļ10/.data/audit/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/maekawanaoto/AIťĖčÁôļ10/.data/audit/launchd.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/maekawanaoto/AIťĖčÁôļ10</string>
</dict>
</plist>
